# vps-bot.service
[Unit]
Description=DarkNodes VPS Manager Discord Bot
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Mybot
ExecStart=/usr/bin/python3 /root/Mybot/bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
